# systemd/overseer-telegram.service
[Unit]
Description=Overseer AI Agent - Telegram Bot
After=network.target overseer.service
PartOf=overseer.service

[Service]
Type=simple
WorkingDirectory=/opt/overseer
ExecStart=/usr/bin/npx tsx /opt/overseer/src/bot/index.ts
Restart=on-failure
RestartSec=15
Environment=NODE_ENV=production
EnvironmentFile=/opt/overseer/.env

[Install]
WantedBy=multi-user.target
